# 3-Main-SoC-Realtek-RTL8196E/33-Rootfs/skeleton/etc/inittab
# /etc/inittab - init configuration for busybox init
#
# Minimal inittab in squashfs rootfs - boots system and mounts /userdata
# Full configuration is in /userdata/etc/inittab (loaded after mount)
#
::sysinit:/etc/init.d/rcS
ttyS0::askfirst:/bin/login
